# tc2/device-register/device-register-overwrite.service
[Service]
ExecStart=
ExecStart=/usr/bin/device-register --reboot --prefix=tc2